# www/framework/Fs/Tests/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:vsftpd]
command=vsftpd
